# roles/docker/templates/containerd.service.j2
[Service]
ExecStartPre=
ExecStartPre=-/sbin/modprobe overlay
