# roles/docker/templates/docker.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H unix:// -s overlay
